# root/etc/s6-overlay/s6-rc.d/init-usermod/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mv /app/data /tmp/data

if [[ ! -f "/data/bot/bot.yml" ]]; then
    echo "[NOTICE] Bot data not found in /data/yaml folder, creating this folder and copying data"
    mv /tmp/data/ /data/bot
fi

ln -s /data/bot /app/data

rm -rf /data/bot/strings
cp -R /tmp/data/strings/ /data/bot/strings

chown -R app:app \
    /app \
    /config \
    /data \
    /app/data
