# root/etc/s6-overlay/s6-rc.d/init-usermod/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

chown -R app:app \
    /app \
    /config \
    /data
# #!/command/execlineb -P
# bash /etc/s6-overlay/s6-rc.d/init-usermod/fixperms.sh